# misc/supervisor.ini
[program:feed-api]
process_name = feed-api
command = /sites/feed-api/feed-api server
directory = /sites/feed-api/
user = archive
autostart = true
autorestart = true
redirect_stderr = true
stdout_logfile = /sites/feed-api/logs/feed-api.log
